# config/autostart/copy-paste.desktop
[Desktop Entry]
Comment=Enable copy paste from/to host (susah nyarinya :D)
Exec=autocutsel -fork
Terminal=false